# services/confd/confd.service
[Unit]
Description=Confd
After=haproxy.service

[Service]
ExecStart=bin/confd
Restart=always

[Install]
WantedBy=basic.target
