# services/rsync/rsync.target
[Unit]
Description=Rsync Timer Target
StopWhenUnneeded=yes
