# services/torrent/config/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:php-fpm]
command=/usr/sbin/php5-fpm

[program:torrent]
command=/bin/bash -c "/etc/init.d/torrent-rtorrent start"
